# daily-schedule-scraper.timer
[Unit]
Description=Run Daily Schedule Scraper hourly

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target 